# deploy/systemd/trader-once.timer
[Unit]
Description=Run trader-once every 15 minutes + 15 seconds

[Timer]
OnCalendar=*:0/15:15
AccuracySec=1s
Persistent=true
Unit=trader-once.service

[Install]
WantedBy=timers.target

