# systemd/trader-once.service
[Unit]
Description=Crypto Strategy Realtime Once (venv)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/crypto_strategy_project
ExecStartPre=/usr/bin/find /opt/crypto_strategy_project -name __pycache__ -type d -exec rm -rf {} +
ExecStart=/opt/crypto_strategy_project/.venv/bin/python /opt/crypto_strategy_project/scripts/realtime_loop.py \
  --cfg /opt/crypto_strategy_project/csp/configs/strategy.yaml --delay-sec 15 --once
Environment=PYTHONUNBUFFERED=1
Restart=no

[Install]
WantedBy=multi-user.target
