# systemd/trader.timer
[Unit]
Description=Run trader.service every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true
AccuracySec=1s
Unit=trader.service

[Install]
WantedBy=timers.target
